# roles/patchwave/templates/patchwave.service.j2
[Unit]
Description=PatchWave Auto-Patching Service
After=network.target

[Service]
Type=oneshot
ExecStart={{ patchwave_script_dir }}/patchwave-patch.sh
StandardOutput=journal
StandardError=journal
